# rootupd.cron
54 21 * * * mv /etc/bind/db.root /etc/bind.db.root_old && wget -q -0 /etc/bind/db.root https://www.internic.net/zones/named.root && service bind9 reload
